# mailpit
#!/sbin/openrc-run
name="Monitor Directory"
description="Monitor a directory"
command="/usr/local/bin/mailpit -d /var/lib/mailpit/mailpit.db -m 1000 --smtp-auth-allow-insecure"
command_background=true
pidfile="/run/monitor_directory.pid"
