# mailpit
#!/sbin/openrc-run
name="Mailpit Service"
description="Mailpit Service"
command="/usr/local/bin/mailpit.sh"
command_background=true
pidfile="/run/mailpit.pid"
